# api/docker/app.conf
[program:api]
command=/bin/sh -c "uvicorn main:app --host 0.0.0.0 --port ${PORT:-8000} --no-access-log"
directory=/app/api
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=10
process_name=%(program_name)s
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
priority=30
; api wraps uvicorn under /bin/sh, and uvicorn itself forks workers. Without
; group signaling, supervisorctl restart races: the new uvicorn tries to
; bind :8000 before the old workers have released it ("address already in
; use"). stopasgroup+killasgroup signal the whole process group; stopwait
; gives uvicorn time for graceful shutdown before SIGKILL.
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true
environment=PATH="/app/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"