# api/docker/cloudflared.conf
[program:cloudflared]
command=/usr/local/bin/cloudflared-start.sh
autostart=false
autorestart=true
startretries=5
startsecs=5
process_name=%(program_name)s
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
priority=60
; Tunnel daemon graceful-stops on SIGTERM. Group signaling matches the
; pattern used for api/ui/workers.
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
