# api/docker/postgres.conf
[program:postgres]
command=/usr/lib/postgresql/18/bin/postgres -D %(ENV_SHS_PG_DATA)s
user=postgres
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=5
process_name=%(program_name)s
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
priority=10
; Postgres uses SIGINT for "fast shutdown" - disconnects clients and exits
; promptly. SIGTERM ("smart shutdown") would wait for clients to disconnect,
; which can hang restart for minutes. Postgres manages its own children so
; group signaling is unnecessary.
stopsignal=INT
stopwaitsecs=10
