# api/docker/ui.conf
[program:ui]
command=/usr/local/bin/ui-start.sh
directory=/app/ui
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=5
process_name=%(program_name)s
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
priority=40
; Group signaling matches api/workers - restart should always release :3000
; cleanly before re-bind. Node's signal handling is fast but adding the
; group settings is consistent and harmless.
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
; NEXT_PUBLIC_* are NOT set here - ui-start.sh renders /app/ui/public/__env.js
; from /workspace/.env at every restart so the browser picks up runtime values
; (CF tunnel domain, etc.). ui/shared/lib/config.ts prefers window.__ENV over
; the build-time NEXT_PUBLIC_*.
environment=NODE_ENV="production",HOSTNAME="0.0.0.0",PORT="3000"
