# api/docker/worker-general.conf
[program:worker-general]
command=python -m shared.worker
directory=/app/worker
autostart=true
autorestart=true
startretries=5
numprocs=%(ENV_SHS_GENERAL_WORKERS)s
startsecs=5
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
priority=50
; Workers have a graceful-shutdown handler that deregisters from the API
; on SIGTERM. Group signaling ensures any subprocess (e.g. ffmpeg in a
; future video worker) goes down with the parent.
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true
environment=SHS_WORKER_TYPE="general",SHS_WORKSPACE_ROOT="/workspace",SHS_API_BASE_URL="http://localhost:8000",PATH="/app/api/.venv/bin:/usr/local/bin:/usr/bin:/bin",PYTHONPATH="/app/worker:/app"
